# bot.conf
[program:bot]
directory=/home/bot/code/vpn_dan_bot
command=/home/bot/code/vpn_dan_bot/_main.py
user=bot
autostart=true
autorestart=true