# server/bot_subserver.conf
[program:bot_subserver]
directory=/home/bot/code/vpn_dan_bot
command=./_serv.py
user=bot
autostart=true
autorestart=true